# prusa_status.service
[Unit]

Description=printer status for prusa

[Service]
ExecStart=pwd
Restart=always

[Install]
WantedBy = multi-user.target
